# deploy/linux/systemd/vllm.service
[Unit]
Description=vLLM AWQ inference service
After=network-online.target
Wants=network-online.target
Conflicts=llama-cpp.service

[Service]
Type=simple
EnvironmentFile=%h/.config/bordercore/vllm.env
ExecStartPre=-/usr/bin/docker rm -f bordercore-vllm
ExecStartPre=/usr/bin/mkdir -p /home/jerrell/.cache/vllm
ExecStart=/usr/bin/docker run --rm --name bordercore-vllm --gpus all --ipc=host -p 127.0.0.1:8001:8000 -v /home/jerrell/models:/models:ro -v /home/jerrell/.cache/vllm:/root/.cache/vllm vllm/vllm-openai@sha256:fc56161ee42a011aeee78b65d0a81b6683c7d04402fd40503d14d4d6c98f07cb ${VLLM_MODEL_PATH} --served-model-name ${VLLM_SERVED_MODEL_NAME} --gpu-memory-utilization ${VLLM_GPU_MEMORY_UTILIZATION} --max-model-len ${VLLM_MAX_MODEL_LEN} --max-num-seqs ${VLLM_MAX_NUM_SEQS} --hf-overrides ${VLLM_HF_OVERRIDES} --generation-config vllm
ExecStop=/usr/bin/docker stop -t 20 bordercore-vllm
ExecStopPost=-/usr/bin/docker rm -f bordercore-vllm
Restart=on-failure
RestartSec=10
TimeoutStartSec=300
TimeoutStopSec=30

[Install]
WantedBy=default.target
